# Lesson06/service_scan/systemd/system/myscan.timer
[Unit]
Description=Execute script test scan log...

[Timer]

Unit=myscan.service
OnCalendar=*-*-* *:*:10,40
[Install]
WantedBy=multi-user.target
